# deploy/proxmox-dashboard.service
[Unit]
Description=Proxmox Dashboard
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/proxmox-dashboard
EnvironmentFile=/home/ubuntu/proxmox-dashboard/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
